# scripts/backup/wordpress/database_dump.timer
[Unit]
Description=Periodical database dumps
Requires=database_dump.service

[Timer]
Unit=database_dump.service
OnCalendar=Daily

